# Final-Project/deployment/motion-detector.service
[Unit]
Description=Gang Motion Detection Core Service
# Memastikan service berjalan hanya setelah network dan filesystem siap
After=multi-user.target

[Service]
# Security: Menjalankan service sebagai user biasa (bukan root)
User=pi

# Mendefinisikan root directory proyek agar import module berjalan lancar
WorkingDirectory=/home/pi/camera_project

# CRITICAL: Menggunakan Python Interpreter dari Virtual Environment (venv)
# Ini memastikan dependensi terisolasi dan tidak konflik dengan system python
ExecStart=/home/pi/camera_project/bin/python3 src/motion_detector.py

# Reliability: Service akan otomatis restart jika crash/error
Restart=always

# Memberikan jeda 10 detik sebelum restart untuk mencegah 'restart loop' berlebihan
RestartSec=10

[Install]
# Mengaktifkan service agar otomatis jalan saat booting (runlevel 3/4/5)
WantedBy=multi-user.target